# marketing/carousel-growth-engine/com.grubpos.carousel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grubpos.carousel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sam/Grub24/agency-agents/marketing/carousel-growth-engine/local-daily.sh</string>
    </array>

    <!-- Daily at 18:00 local time (6pm) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Do not run when the plist is loaded; only on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/sam/Grub24/agency-agents/marketing/carousel-growth-engine/state/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sam/Grub24/agency-agents/marketing/carousel-growth-engine/state/launchd.err.log</string>
</dict>
</plist>
